# etc/github/templates/etc/profile.d/orchestrator-client.sh.erb
export ORCHESTRATOR_API='<%= orchestrator_api_nodes.join(" ") %>'