# data/github/current/config/supervisord.conf
# The vtcombo process is set to auto restart. In case mysql is down it will be in a restart loop
# with the state "RUNNING" in supervisorctl.
[program:vtcombo]
redirect_stderr=true
command=%(here)s/../script/vtcombo-up
user=vitess
group=vitess
