# etc/github/templates/lib/systemd/system/cluster-rebalance.timer.erb
[Unit]
Description=cluster-balance daemon

[Timer]
OnCalendar=<%= cluster_rebalance_schedule %>

[Install]
WantedBy=github-enterprise.target

