# etc/github/templates/lib/systemd/system/cluster-rebalance.service.erb
[Unit]
Description=ghe-cluster-balance daemon
Requires=nomad.service
After=nomad.service

[Service]
Restart=on-failure
ExecStart=ghe-cluster-balance rebalance --yes <% if cluster_rebalance_apps? %> --job-name <%= cluster_rebalance_apps %> <% end %><% if cluster_rebalance_workers? %> --workers <%= cluster_rebalance_workers %><% end %><% if cluster_rebalance_timeout? %> --timeout <%= cluster_rebalance_timeout %><% end %>
Type=oneshot

[Install]
WantedBy=multi-user.target
