# usr/lib/systemd/system/elasticsearch-upgrade.service
[Unit]
Description=ElasticSearch indices upgrade

[Service]
Type=oneshot
ExecStart=/usr/local/share/enterprise/ghe-es-upgrade
