# usr/lib/systemd/system/ghe-user-disk.path
[Unit]
Description=Monitor availability of user disk
Requires=data-user.mount
After=data-user.mount

[Path]
PathExists=/data/user/version
Unit=ghe-user-disk.service
