# usr/lib/systemd/system/ghe-health-monitor.timer
[Unit]
Description=GHE monitoring binary

[Timer]
OnCalendar=minutely

[Install]
WantedBy=github-enterprise.target
