# usr/lib/systemd/system/ghes-manage-agent.service
[Unit]
Description=GHES Manage Agent
ConditionPathExists=/data/user/version
ConditionPathExists=/data/ghes-manage-agent/shared
ConditionPathIsMountPoint=/data/user
BindsTo=ghe-user-disk.service
Wants=ghe-user-disk.service
After=ghe-user-disk.service

[Service]
User=ghes-manage-agent
Type=simple
Restart=always
RestartSec=5s
RestartPreventExitStatus=126
KillMode=control-group
KillSignal=SIGTERM
ExecReload=/bin/kill -USR2 $MAINPID
TimeoutStartSec=180
TimeoutStopSec=600
SendSIGKILL=no

WorkingDirectory=/data/ghes-manage-agent/shared
LimitCORE=524288000

Environment=CHECK_SERVICE=no
SyslogIdentifier=ghes-manage-agent
ExecStart=/usr/local/share/enterprise/ghe-systemd-wrapper ghes-manage-agent /usr/local/share/enterprise/ghes-manage \
  --agent

[Install]
WantedBy=multi-user.target
