# usr/lib/systemd/system/github-enterprise-ha.target
[Unit]
Description=GitHub Enterprise cluster-based HA
