# usr/lib/systemd/system/github-enterprise.target
[Unit]
Description=GitHub Enterprise
Wants=ghe-user-disk.service
Wants=enterprise-manage.service

[Install]
WantedBy=multi-user.target
